# contrib/systemd/btrbk.service.in
[Unit]
Description=btrbk backup
After=network.target

[Service]
Type=oneshot
ExecStart=@BINDIR@/btrbk run
