# contrib/systemd/btrbk-network.service.in
[Unit]
Description=btrbk backup
Documentation=man:btrbk(1)
After=network-online.target

[Service]
Type=oneshot
ExecStart=@BINDIR@/btrbk run
