# contrib/systemd/btrbk.service.in
[Unit]
Description=btrbk backup

[Service]
Type=oneshot
ExecStart=@BINDIR@/btrbk run
