# contrib/systemd/btrbk.service.in
[Unit]
Description=btrbk backup
Documentation=man:btrbk(1)
ConditionPathExists=|@CONFDIR@/btrbk/btrbk.conf
ConditionPathExists=|@CONFDIR@/btrbk.conf

[Service]
Type=oneshot
ExecStart=@BINDIR@/btrbk run
