# contrib/systemd/btrbk.timer.in
[Unit]
Description=btrbk daily backup

[Timer]
OnCalendar=hourly
AccuracySec=10min
Persistent=true

[Install]
WantedBy=timers.target
