# contrib/systemd/btrbk.service.in
[Unit]
Description=btrbk backup
Documentation=man:btrbk(1)

[Service]
Type=oneshot
ExecStart=@BINDIR@/btrbk run
