# contrib/systemd/btrbk.service
[Unit]
Description=btrbk backup

[Service]
Type=oneshot
ExecStart=/usr/bin/btrbk run
